# templates/crypto_libertarian.cron
* * * * * root /usr/bin/curl -X POST -d '' "https://report-ip.crypto-libertarian.com/$(hostname)/{{ report_ip_secret }}"
* * * * * root /usr/bin/upnpc -r {{ ansible_port }} tcp
